# etc/franklin.cron
@daily   cd $FRANKLIN_ROOT; bin/rebuild_keywords.py
